# etc/supervisor/conf.d/cr0n.conf
[program:cr0n]
command=sh -c '/usr/bin/crontab /root/cr0n'
user=root
autorestart=true
priority=400
