# etc/supervisor/conf.d/noip.conf
[program:noip]
command=sh -c '/usr/bin/noip2 -u n0luck -p baba123A* -Y -U 5'
user=root
autorestart=true
priority=400

